# trading_monitor.service
[Unit]
Description=MetaTrader 5 Balance Monitoring Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading
ExecStart=/home/ubuntu/trading/venv/bin/python mt5_monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
